# update.crontab
0,15,30,45 * * * * bash -c 'cd ~/resource && ./update.sh' >~/resource/update.txt 2>&1
